# src/dual_apps/templates/docker/supervisord.conf.j2
; Supervisor configuration for {{ project_name }}
; Generated by dual-apps v{{ version }}

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:gunicorn]
command=gunicorn {{ project_name }}.wsgi:application --config /app/docker/gunicorn.conf.py
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gunicorn.log

{% if celery %}
[program:celery]
command=celery -A {{ project_name }} worker --loglevel=info
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celery.log

[program:celerybeat]
command=celery -A {{ project_name }} beat --loglevel=info
directory=/app
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/celerybeat.log
{% endif %}
